# etc/profile.d/console-login-helper-messages-profile.sh
../../usr/share/console-login-helper-messages/profile.sh