# etc/systemd/system/kubelet.service.d/20-nodenet.conf
[Service]
Environment="KUBELET_NODE_IP=10.19.115.9" "KUBELET_NODE_IPS=10.19.115.9,2620:52:0:1372:250:56ff:fe89:6230"
