# etc/systemd/system/rpm-ostreed.service.d/mco-controlplane-nice.conf
# See https://github.com/openshift/machine-config-operator/issues/1897
[Service]
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=6
